# xebikart-car-setup/roles/xebikart-car-app/templates/xebikart-car-app.service
[Unit]
Description={{ app_name }}

[Service]
User=pi
ExecStart={{ user_directory }}/env/bin/python {{ user_directory }}/{{ app_name }}/manage.py drive

[Install]
WantedBy=multi-user.target